# luci/trunk/package/luci/applications/luci-transmission/root/etc/init.d/trextra
#!/bin/sh /etc/rc.common

START=99

start(){
	/usr/bin/trextra start
}

stop(){
	/usr/bin/trextra stop
}
restart(){
stop
start
}

