# luci/trunk/themes/XEyE/root/etc/uci-defaults/luci-theme-xeye
#!/bin/sh
uci batch <<-EOF
	uci set luci.themes.LuciXEyE=/luci-static/LuciXEyE
	commit luci
EOF

