# luci/trunk/themes/xeye/root/etc/uci-defaults/luci-theme-xeye
#!/bin/sh
uci batch <<-EOF
	uci set luci.themes.LuciXEyE=/luci-static/LuciXEyE
	uci set luci.main.mediaurlbase=/luci-static/LuciXEyE
	uci commit
EOF

